# scripts/virtual-office-buddy.service
# systemd service file for Virtual Office Buddy (self-hosted, no Docker)
#
# Install:
#   sudo cp scripts/virtual-office-buddy.service /etc/systemd/system/virtual-office-buddy.service
#   # Edit User= and WorkingDirectory= below to match your setup
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now virtual-office-buddy
#
# Logs:
#   sudo journalctl -u virtual-office-buddy -f

[Unit]
Description=Virtual Office Buddy
After=network.target

[Service]
Type=simple
User=root
Environment="HOME=/root"
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
WorkingDirectory=/root/virtual-office-buddy

# Load environment from .env if present
EnvironmentFile=-/root/virtual-office-buddy/.env

ExecStart=/usr/bin/npm run preview -- --port 3000 --host 0.0.0.0
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
